# deploy/face-access.service
[Unit]
Description=RP3 Face Access Control System
Documentation=https://github.com/your-repo/rp3-face-access
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/hseadmin/RPBluetooth
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/hseadmin/RPBluetooth/src

# Set static IP before starting (ignore if already set)
ExecStartPre=-/sbin/ip addr add 192.168.1.22/24 dev eth0

# Main executable (runs as root)
ExecStart=/home/hseadmin/RPBluetooth/venv/bin/python -u src/main.py --config config/production.yaml

# Restart policy - always restart on failure
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Watchdog - service must ping within 60 seconds or gets restarted
WatchdogSec=60

# Stop timeout
TimeoutStopSec=10

# Resource limits
MemoryMax=512M
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=face-access

[Install]
WantedBy=multi-user.target
